# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fastapi]
directory=/fastapi
command=/bin/bash -c "/wait-for-port.sh 127.0.0.1 5678 && exec uvicorn main:app --host 0.0.0.0 --port 9374"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
redirect_stderr=true
priority=10

[program:n8n]
directory=/home/node
command=n8n start
autostart=true
autorestart=true
environment=DB_TYPE=sqlite,DB_SQLITE_DATABASE=/home/node/.n8n/database.sqlite
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
redirect_stderr=true
priority=20

[program:nginx]
command=/bin/bash -c "/wait-for-port.sh 127.0.0.1 9374 && /wait-for-port.sh 127.0.0.1 5678 && exec nginx -g 'daemon off;'"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
redirect_stderr=true
priority=30
